# raspi-eink-dashboard.service
[Unit]
Description=E-ink Dashboard Web Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=justin
WorkingDirectory=/home/justin/raspi-eink-dashboard
ExecStart=/usr/bin/python3 /home/justin/raspi-eink-dashboard/server.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
